# docker/s6-overlay/s6-rc.d/addon-migration/up
#!/command/execlineb -P

# Run addon migrations after Coolify's main migrations
with-contenv
cd /var/www/html
foreground {
    php
    artisan
    migrate
    --path=vendor/amirhmoradi/coolify-enhanced/database/migrations
    --force
}
