# docker/s6-overlay/s6-rc.d/addon-migration/up
#!/command/execlineb -P

# Run addon migrations after Coolify's main migrations
with-contenv
cd /var/www/html
foreground {
    php
    artisan
    migrate
    --path=vendor/amirhmoradi/coolify-enhanced/database/migrations
    --force
}

# Create custom templates cache directory (www-data needs write access)
foreground { mkdir -p /data/coolify/custom-templates }
foreground { chown www-data:www-data /data/coolify/custom-templates }
